# aux-scope-example.service
[Service]
Type=notify
KillMode=mixed
ExecStart=/usr/local/bin/aux-scope-example
